# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=1

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app
autorestart=true
priority=2
